# TwinningAgent/camera_probe.service
[Unit]
Description=Camera health probe
After=network-online.target

[Service]
ExecStart=/usr/bin/python3 /usr/local/bin/camera_probe.py \
  --ip 192.168.254.3 \
  --rtsp rtsp://192.168.20.2:8554/proxied \
  --broker 192.168.20.2 \
  --interval 10 --fail 3 --recover 3
Restart=always
User=root

[Install]
WantedBy=multi-user.target
